# gooroom/gooroom-3.x/gooroom-feedback-server/debian/gooroom-feedback-server.postinst
#!/bin/sh
set -e

/bin/systemctl daemon-reload || true
/bin/systemctl enable gooroom-feedback-server || true
/bin/systemctl start gooroom-feedback-server || true

exit 0
